# scripts/magi-archive-uploader.service
[Unit]
Description=Magi Archive Uploader (one-shot pass)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/magi-archive-uploader
EnvironmentFile=/home/ubuntu/magi-archive-uploader/.env
ExecStart=/usr/bin/node /home/ubuntu/magi-archive-uploader/dist/index.js
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/.local/share/magi-archive-uploader /home/ubuntu/magi-assistant-foundry/data/video /home/ubuntu/magi-assistant-discord/data/sessions
PrivateTmp=true

[Install]
WantedBy=multi-user.target
